# systemd/jinwang-jarvis-hermes-health.service
[Unit]
Description=Check Hermes gateway + Jarvis cron health and alert Discord
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/jinwang/workspace/jinwang-jarvis
Environment=PATH=/home/jinwang/.local/bin:/home/jinwang/.local/bin:/home/jinwang/.hermes/hermes-agent/venv/bin:/home/jinwang/.hermes/hermes-agent/node_modules/.bin:/home/jinwang/.local/share/fnm/node-versions/v24.14.1/installation/bin:/home/jinwang/.local/bin:/home/jinwang/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin
Environment=HERMES_HOME=/home/jinwang/.hermes
Environment=JARVIS_HEALTH_DISCORD_CHANNEL=1496014213276241922
EnvironmentFile=-/home/jinwang/.hermes/.env
ExecStart=/bin/bash -lc 'cd /home/jinwang/workspace/jinwang-jarvis && PYTHONPATH=src python3 -m jinwang_jarvis.cli hermes-health-check --config config/pipeline.local.yaml --discord-alert --restart --stale-minutes 15'
