# systemd/zeus-os-hermes-health.service
[Unit]
Description=Check Hermes gateway + ZeusOS cron health and alert Discord
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/jinwang/workspace/zeus-os
Environment=PATH=/home/jinwang/.hermes/hermes-agent/venv/bin:/home/jinwang/.hermes/hermes-agent/node_modules/.bin:/home/jinwang/.local/share/fnm/node-versions/v24.14.1/installation/bin:/home/jinwang/.local/bin:/home/jinwang/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HERMES_HOME=/home/jinwang/.hermes
Environment=ZEUSOS_HEALTH_DISCORD_CHANNEL=1496014213276241922
EnvironmentFile=-/home/jinwang/.hermes/.env
ExecStart=/bin/bash -lc 'cd /home/jinwang/workspace/zeus-os && PYTHONPATH=src python3 -m zeus_os.cli hermes-health-check --config config/pipeline.local.yaml --discord-alert --restart --stale-minutes 15 --readiness-timeout-seconds 45'
