# orchestration/2026-05-07-zeusos-full-pivot/backups/systemd-stage3-20260507T085927Z/jinwang-jarvis-hermes-health.timer
[Unit]
Description=Run Hermes/Jarvis health check every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true
Unit=jinwang-jarvis-hermes-health.service

[Install]
WantedBy=timers.target
