# orchestration/2026-05-07-zeusos-full-pivot/backups/systemd-stage3-20260507T085927Z/jinwang-jarvis-cycle.service
[Unit]
Description=Jinwang Jarvis pipeline polling cycle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/tmp/pytest-of-jinwang/pytest-846/test_cli_run_cycle_and_install0
Environment=PATH=/home/jinwang/.local/bin:/home/jinwang/.hermes/hermes-agent/venv/bin:/home/jinwang/.hermes/hermes-agent/node_modules/.bin:/home/jinwang/.local/share/fnm/node-versions/v24.14.1/installation/bin:/home/jinwang/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin
ExecStart=/bin/bash -lc 'cd /tmp/pytest-of-jinwang/pytest-846/test_cli_run_cycle_and_install0 && PYTHONPATH=src python3 -m jinwang_jarvis.cli run-cycle --config pipeline.yaml'
