# orchestration/2026-05-07-zeusos-full-pivot/backups/systemd-stage3-20260507T085927Z/zeus-os-weekly-review.timer
[Unit]
Description=Run ZeusOS weekly review on Sunday evening

[Timer]
OnCalendar=Sun *-*-* 20:00:00
Persistent=true
Unit=zeus-os-weekly-review.service

[Install]
WantedBy=timers.target
