# orchestration/2026-05-07-zeusos-full-pivot/backups/systemd-stage3-20260507T085927Z/zeus-os-cycle.timer
[Unit]
Description=Run ZeusOS pipeline every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
Persistent=true
Unit=zeus-os-cycle.service

[Install]
WantedBy=timers.target
